# root/etc/init.d/zabbix-agentd-config
#!/bin/sh /etc/rc.common

START=95
STOP=01

USE_PROCD=1

start_service() {
	config_load 'zabbix-agentd'
	
	# 创建配置文件目录
	mkdir -p /etc/zabbix/zabbix_agentd.d/
	
	# 生成主配置文件
	config_get_bool enable general enable 0
	if [ $enable -eq 1 ]; then
		cat > /etc/zabbix/zabbix_agentd.conf << EOF
# Zabbix Agentd Configuration
# Generated by LuCI

PidFile=/var/run/zabbix/zabbix_agentd.pid
LogFileSize=0

EOF
		
		# 添加基本配置项
		config_get server general server "127.0.0.1"
		echo "Server=$server" >> /etc/zabbix/zabbix_agentd.conf
		
		config_get server_active general server_active "127.0.0.1:10051"
		echo "ServerActive=$server_active" >> /etc/zabbix/zabbix_agentd.conf
		
		config_get hostname general hostname
		if [ -n "$hostname" ]; then
			echo "Hostname=$hostname" >> /etc/zabbix/zabbix_agentd.conf
		else
			echo "Hostname=`hostname`" >> /etc/zabbix/zabbix_agentd.conf
		fi
		
		config_get listen_ip general listen_ip "0.0.0.0"
		echo "ListenIP=$listen_ip" >> /etc/zabbix/zabbix_agentd.conf
		
		config_get listen_port general listen_port "10050"
		echo "ListenPort=$listen_port" >> /etc/zabbix/zabbix_agentd.conf
		
		config_get_bool enable_remote_commands general enable_remote_commands 0
		echo "EnableRemoteCommands=$enable_remote_commands" >> /etc/zabbix/zabbix_agentd.conf
		
		config_get log_level general log_level "3"
		echo "LogLevel=$log_level" >> /etc/zabbix/zabbix_agentd.conf
		
		config_get log_file general log_file "/var/log/zabbix-agentd.log"
		echo "LogFile=$log_file" >> /etc/zabbix/zabbix_agentd.conf
		
		config_get timeout general timeout "30"
		echo "Timeout=$timeout" >> /etc/zabbix/zabbix_agentd.conf
		
		config_get include_dir general include_dir "/etc/zabbix/zabbix_agentd.d/"
		echo "Include=$include_dir*.conf" >> /etc/zabbix/zabbix_agentd.conf
		
		# 添加用户自定义参数
		config_foreach add_userparameter userparameter
	fi
}

add_userparameter() {
	config_get name "$1" name
	config_get command "$1" command
	
	if [ -n "$name" ] && [ -n "$command" ]; then
		echo "UserParameter=$name,$command" >> /etc/zabbix/zabbix_agentd.conf
	fi
}

reload_service() {
	start_service
	/etc/init.d/zabbix-agentd restart
}